# imagebuilder/assets/cloud-init/etc/update-motd.d/10-help-text
#!/bin/sh
printf "\n"
printf " * NTP Sources:    chronyc sources\n"
printf " * GPS Status:     gpsmon\n"